# scripts/com.colinvinson.rowan-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.colinvinson.rowan-backup</string>

  <!-- Run the weekly Supabase backup script. -->
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/colinvinson/rowan-dashboard/node_modules/.bin/tsx</string>
    <string>/Users/colinvinson/rowan-dashboard/scripts/backup.ts</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/colinvinson/rowan-dashboard</string>

  <!-- Fire every Sunday at 03:00 local time. If the Mac is asleep, launchd
       runs it at the next wake. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>0</integer>
    <key>Hour</key><integer>3</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/colinvinson/rowan-dashboard/.rowan-backup.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/colinvinson/rowan-dashboard/.rowan-backup.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>Nice</key><integer>10</integer>
</dict>
</plist>
